# orchestrator/com.octopoid.scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.octopoid.scheduler</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>-m</string>
        <string>orchestrator.scheduler</string>
        <string>--once</string>
        <string>--debug</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/maxwilliams/dev/octopoid</string>

    <key>StartInterval</key>
    <integer>60</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/maxwilliams/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/maxwilliams/dev/octopoid</string>
        <key>HOME</key>
        <string>/Users/maxwilliams</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/maxwilliams/dev/octopoid/.octopoid/runtime/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/maxwilliams/dev/octopoid/.octopoid/runtime/logs/launchd-stderr.log</string>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
